# deployments/docker/supervisord.conf
; supervisord configuration for the radius-manager-api Docker dev image.
;
; The image runs supervisord as PID 1 (entrypoint.sh execs it). supervisord
; in turn launches:
;   - rm-api          : the radius-manager-api HTTP server (`serve`).
;   - freeradius      : the FreeRADIUS daemon in foreground/debug mode.
;
; Per-instance freeradius-api programs are added at runtime by the Go
; SupervisordSystemctl backend, which drops [program:<name>-api] sections
; into /etc/supervisor/conf.d/<name>-api.conf and shells `supervisorctl
; reread && supervisorctl update`.

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[include]
files = /etc/supervisor/conf.d/*.conf

; ---- Long-lived programs the image always runs ----

[program:freeradius]
; -X = single-process foreground debug mode; supervisord captures stdout+stderr.
; -f = stay foreground (no daemonize). Without -f, freeradius forks and
;      supervisord would mistake the parent's exit for a crash.
command=/usr/sbin/freeradius -X -f
user=root
autorestart=true
autostart=true
startsecs=2
priority=20
stdout_logfile=/var/log/supervisor/freeradius.out.log
stderr_logfile=/var/log/supervisor/freeradius.err.log

[program:rm-api]
command=/usr/local/bin/radius-manager-api serve
user=root
autorestart=true
autostart=true
startsecs=2
priority=30
; rm-api needs the same env vars docker-compose provides; they are
; inherited from the supervisord parent process (PID 1).
environment=HOME="/root"
stdout_logfile=/var/log/supervisor/rm-api.out.log
stderr_logfile=/var/log/supervisor/rm-api.err.log
